# resources/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=/entrypoint.sh mysqld

[program:sshd]
command=/usr/sbin/sshd -D
